# service/s6-svscan-log/run
#!/usr/bin/execlineb -P

if { mkdir -p /tmp/uncaught-logs }
redirfd -w 2 /dev/console
redirfd -rnb 0 fifo
exec -c
s6-log -bpt /tmp/uncaught-logs

