# service/.s6-svscan/finish
#!/usr/bin/execlineb -S0

cd /
fdclose 0
redirfd -w 1 /dev/console
fdmove -c 2 1

foreground { echo "Syncing disks." }
foreground { sync }
foreground { echo "Sending all processes the TERM signal." }
foreground { killall -v -r .* }
foreground { sleep 5 }
foreground { echo "Sending all processes the KILL signal." }
foreground { killall -9 -v -r .* }

wait { }

foreground { echo "Halting" }

